# roles/member/files/rc.local
#!/bin/sh -e
mount /dev/sda1 /opt/crust/data/files/
